# app/gunicorn.conf
[program:gunicorn]
command=gunicorn -c /etc/gunicorn/gunicorn.conf.py app.main:app
directory=/var/www/
autostart=true
autorestart=true
redirect_stderr=true